# init/uranus-web.service
[Unit]
Description=Uranus Web
After=hackernel.service

[Service]
ExecStart=/usr/sbin/uranus-web

[Install]
WantedBy=multi-user.target
