# init/uranus-web.service
[Unit]
Description=Uranus Web
After=hackernel.service

[Service]
ExecStart=/usr/sbin/uranus-web
Restart=always
RestartSec=30
Environment="GOMEMLIMIT=50MiB"

[Install]
WantedBy=multi-user.target
